# templates/default/chef-guard.upstart.erb
# $Id: Generated by chef for node: <%= node['hostname'] %>
# --------------------------------------------------------------
# NOTE: This file is controlled by chef templates!
#   Do not edit or change this file but change the following:
#    * template file: <%= @template_file %>
#    * recipe file: <%= @recipe_file %>
# --------------------------------------------------------------

description "Chef-Guard"
start on runlevel [2345]
stop on runlevel [!2345]

env LANG=en_US.UTF-8
export LANG

respawn

exec <%= @basedir %>/chef-guard >> <%= @basedir %>/console.log 2>&1
